# deploy/systemd/codeblack-celery.service
[Unit]
Description=CodeBlack Celery Worker
After=network-online.target redis-server.service codeblack-backend.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/celery -A celery_worker.celery_app worker --loglevel=info --pool=solo
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
