# deploy/systemd/codeblack-bot.service
[Unit]
Description=CodeBlack Discord Bot
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv-bot/bin/python __APP_DIR__/main.py
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
