# deploy/systemd/codeblack-backend.service
[Unit]
Description=CodeBlack Backend API (FastAPI)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__APP_DIR__/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=__APP_DIR__/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 13400
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
